# utils/daily_upload_download/cassandra_upload/daily_cassandra_upload.cron
01 03 * * * /usr/bin/bash /nonsec/sstable-migrator/utils/daily_upload_download/daily_cassandra_upload.sh