# utils/daily_upload_download/daily_certstream_download.cron
01 00 * * * /usr/bin/bash /root/dnsdata/make_domains.sh